# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0

; ── PHP-FPM ──────────────────────────────────────────────────────────────────
[program:php-fpm]
command=/usr/local/sbin/php-fpm --nodaemonize
autostart=true
autorestart=true
priority=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── Nginx ────────────────────────────────────────────────────────────────────
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── Default queue (alerts, webhooks, vendor propagation) ─────────────────────
[program:queue-default]
command=php /var/www/html/artisan queue:work redis --queue=default --sleep=3 --tries=3 --max-time=3600 --memory=128
process_name=%(program_name)s_%(process_num)02d
autostart=true
autorestart=true
numprocs=2
user=www-data
priority=20
stopwaitsecs=60
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── Import queue (Excel / Tally XML — memory-heavy, single worker) ───────────
[program:queue-imports]
command=php /var/www/html/artisan queue:work redis --queue=imports --sleep=5 --tries=2 --max-time=600 --memory=256
autostart=true
autorestart=true
user=www-data
priority=20
stopwaitsecs=120
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── Scheduler (runs every 60 s, dispatches cron-defined commands) ────────────
[program:scheduler]
command=/bin/sh -c 'while true; do php /var/www/html/artisan schedule:run; sleep 60; done'
autostart=true
autorestart=true
user=www-data
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
