# deploy/supervisor.conf
; /etc/supervisor/conf.d/msme-tracker.conf
; Manages Laravel queue workers and the scheduler

; ---------------------------------------------------------------------------
; Default queue — processes all jobs except 'imports'
; Two workers provide parallelism for alerts, webhook processing, etc.
; ---------------------------------------------------------------------------
[program:msme-worker-default]
process_name        = %(program_name)s_%(process_num)02d
command             = /usr/bin/php /var/www/msme-tracker/artisan queue:work redis --queue=default --sleep=3 --tries=3 --max-time=3600 --memory=128
autostart           = true
autorestart         = true
stopasgroup         = true
killasgroup         = true
user                = www-data
numprocs            = 2
redirect_stderr     = true
stdout_logfile      = /var/log/supervisor/msme-worker-default.log
stdout_logfile_maxbytes = 10MB
stdout_logfile_backups  = 3
stopwaitsecs        = 60

; ---------------------------------------------------------------------------
; Import queue — processes Excel import jobs exclusively
; Single worker so memory-heavy imports don't starve other jobs
; ---------------------------------------------------------------------------
[program:msme-worker-imports]
process_name        = %(program_name)s_%(process_num)02d
command             = /usr/bin/php /var/www/msme-tracker/artisan queue:work redis --queue=imports --sleep=5 --tries=2 --max-time=600 --memory=256
autostart           = true
autorestart         = true
stopasgroup         = true
killasgroup         = true
user                = www-data
numprocs            = 1
redirect_stderr     = true
stdout_logfile      = /var/log/supervisor/msme-worker-imports.log
stdout_logfile_maxbytes = 10MB
stdout_logfile_backups  = 3
stopwaitsecs        = 120

; ---------------------------------------------------------------------------
; Scheduler — runs 'artisan schedule:run' every minute
; Laravel's built-in scheduler then dispatches commands based on their cron
; ---------------------------------------------------------------------------
[program:msme-scheduler]
process_name        = %(program_name)s
command             = /bin/bash -c 'while true; do /usr/bin/php /var/www/msme-tracker/artisan schedule:run >> /var/log/supervisor/msme-scheduler.log 2>&1; sleep 60; done'
autostart           = true
autorestart         = true
user                = www-data
numprocs            = 1
redirect_stderr     = true
stdout_logfile      = /var/log/supervisor/msme-scheduler.log
stdout_logfile_maxbytes = 5MB
stdout_logfile_backups  = 2

; ---------------------------------------------------------------------------
; Group shortcut — reload all workers with: supervisorctl restart msme:*
; ---------------------------------------------------------------------------
[group:msme]
programs = msme-worker-default, msme-worker-imports, msme-scheduler
